# particle-io/supervisor.conf
[supervisord]
nodaemon=true

[program:particle_io]
priority=1
command=node /particle-io/saas-app/app.js
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

;[program:mongodb]
;priority=10
;command=mongod
;stdout_logfile=/var/log/supervisor/mongo.log
;stderr_logfile=/var/log/supervisor/mongo.log
;autorestart=true


[program:nginx]
priority=10
command=nginx -g "daemon off;"
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
autorestart=true